# ha_dashboard_toucher/rootfs/etc/s6-overlay/s6-rc.d/ha_dashboard_toucher/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Home Assistant Jinja Compiler
# Runs the jinja compiler
# ==============================================================================
declare -a options
declare config_path

bashio::log.info 'Starting ha_dashboard_toucher...'

# Run the code server
exec /usr/bin/ha_dashboard_toucher