# ha_uglifyjs/rootfs/etc/s6-overlay/s6-rc.d/ha_uglifyjs/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Home Assistant Jinja Compiler
# Runs the jinja compiler
# ==============================================================================
declare -a options
declare config_path

# Declare variables
declare js_files

## Get the 'message' key from the user config options.
js_files=$(bashio::config 'js_files')

bashio::log.info 'Starting ha_uglifyjs...'
bashio::log.info "Watching ${js_files}"

# Run the code server
exec /usr/bin/ha_uglifyjs "${js_files}"